# deploy/systemd/stockquantbot.service.tpl
[Unit]
Description=StockQuantBot FastAPI Service
After=network.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_USER__
WorkingDirectory=__APP_DIR__
Environment=PYTHONUNBUFFERED=1
Environment=TZ=Asia/Shanghai
ExecStart=__APP_DIR__/.venv/bin/uvicorn api:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
